# roles/carbon/templates/upstart-carbon.conf.j2
# This file was generated by Ansible for {{ansible_fqdn}}
# Do NOT modify this file by hand!

description "Graphite-API server"
start on runlevel [2345]
stop on runlevel [!2345]

console log
limit nofile 65536 65536
respawn
setgid {{graphite_api_group}}
setuid {{graphite_api_user}}

exec {{graphite_api_home}}/bin/carbon_cache.py start --config={{graphite_api_home}}/conf/carbon.conf --pidfile={{graphite_api_home}}/storage/carbon-cache.pid

